# infrastructure/incident-controller/deploy/systemd/incident-controller.service
[Unit]
Description=Tiresias Incident Controller
After=network.target tailscaled.service

[Service]
Type=simple
User=cristian
WorkingDirectory=/home/cristian/tiresias-incident-controller
EnvironmentFile=/etc/incident-controller/env
ExecStart=/usr/bin/python3 -m uvicorn src.main:app --host 0.0.0.0 --port 8090 --log-level info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
